# data/ubuntu-location-provider-geoclue2.in
#!/bin/sh
### BEGIN INIT INFO
# Provides:          ubuntu-location-provider-geoclue2
# Required-Start:    $dbus $network
# Required-Stop:     $dbus
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: AGPS provider
# Description:       Provides AGPS for fairphone with Ubuntu Touch
### END INIT INFO

case "$1" in
    start)
        echo "Starting geoclue2 location provider"
        /usr/bin/ubuntu-location-provider-geoclue2 &
        ;;
    stop)
        echo "Stop geoclue2 location provider"
        PID=`ps -ef | grep /usr/bin/ubuntu-location-provider-geoclue2 | grep -v grep | awk '{print $2}'`
	echo "$PID"
        kill -9 $PID
        ;;
esac
exit 0
